# etc/service/nginx/run
#!/bin/sh
set -e
exec /usr/sbin/nginx -g 'daemon off;' >>/var/log/sv-nginx.log 2>&1
